# claudex/launchd/com.tomasmach.cliproxyapi.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tomasmach.cliproxyapi</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/bin/cliproxyapi</string>
    <string>-config</string>
    <string>__HOME__/.config/cliproxyapi/config.yaml</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.local/state/cliproxyapi/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.local/state/cliproxyapi/stderr.log</string>
</dict>
</plist>
